# docker/supervisord-php5-fpm.conf
[program:php5-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
priority=999
username=www-data
autostart=true
autorestart=true
